# deploy/systemd/cortex-city-disk-hygiene.timer
[Unit]
Description=Run Cortex City disk hygiene cleanup hourly

[Timer]
OnCalendar=hourly
RandomizedDelaySec=10m
Persistent=true
Unit=cortex-city-disk-hygiene.service

[Install]
WantedBy=timers.target
