# deploy/systemd/cortex-city-host-metrics.service
[Unit]
Description=Cortex City host metrics logger
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cortex
Group=cortex
WorkingDirectory=/opt/cortex-city/app
Environment=HOST_METRICS_INTERVAL_SECONDS=60
Environment=HOST_METRICS_RETENTION_DAYS=3
Environment=HOST_METRICS_MODE=compact
Environment=HOST_METRICS_DETAIL_EVERY=60
Environment=TMPDIR=/opt/cortex-city/app/tmp
Environment=TMP=/opt/cortex-city/app/tmp
Environment=TEMP=/opt/cortex-city/app/tmp
EnvironmentFile=-/etc/cortex-city/host-metrics.env
ExecStart=/usr/bin/env bash ./scripts/host-metrics-logger.sh
Restart=always
RestartSec=5
TimeoutStopSec=15
KillMode=control-group

[Install]
WantedBy=multi-user.target
